# roles/retro_install/files/etc/profile.d/99-launcher.sh
#!/bin/sh

if [ "$(tty)" == "/dev/tty1" ] ; then
    retroarch
fi

