# roles/retro_install/files/etc/profile.d/99-mymenu.sh
#!/bin/sh

if [ "$(tty)" == "/dev/tty1" ] ; then
    switchres 320 240 60 -s --ini /etc/switchres.ini --launch mymenu
fi

