# roles/retro_install/files/etc/profile.d/99-retroarch.sh
#!/bin/sh

if [ "$(tty)" == "/dev/tty1" ] ; then
    #retroarch
    #startx
    echo "retroarch for kms, startx for xorg"
fi

