# roles/retro_install/files/etc/systemd/system/switchresconf.service
[Unit]
Description=Configure switchres

[Service]
Type=oneshot
ExecStart=/usr/local/bin/switchresconf
